# graphwiz-live-trading.service
[Unit]
Description=GraphWiz Live Trading - Real Money Trading with GoEmotions
Documentation=https://github.com/graphwiz-trader/graphwiz-trader
After=network-online.target docker.service
Wants=docker.service
Requires=network-online.target

[Service]
Type=simple
User=weiss
Group=weiss
WorkingDirectory=/opt/git/graphwiz-trader

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/opt/git/graphwiz-trader/.env

# Pre-start: Build and prepare
ExecStartPre=/usr/bin/docker-compose build
ExecStartPre=-/usr/bin/docker-compose down

# Main: Start full stack with Docker Compose
ExecStart=/usr/bin/docker-compose up -d

# Stop: Graceful shutdown
ExecStop=/usr/bin/docker-compose down

# Reload: Restart services (graceful)
ExecReload=/usr/bin/docker-compose restart trading



# Restart policy (disabled for Type=simple with docker-compose)
# Restart is handled by Docker Compose restart policies in docker-compose.yml

# Logging
StandardOutput=append:/opt/git/graphwiz-trader/logs/live_trading/systemd_service.log
StandardError=append:/opt/git/graphwiz-trader/logs/live_trading/systemd_service_error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

# Timeout for graceful shutdown
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
